# crontab
*/15 * * * * domain-connect-dyndns --config /config.json --all update
